# infra/ansible/roles/app/templates/bestbooks-api.service.j2
# Managed by Ansible (roles/app)
[Unit]
Description=Best Books Guide API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}/current
# Two files, two owners: .env is Ansible's (config), .env.release is the deploy
# playbook's (APP_VERSION). Later files win, and `-` tolerates a missing one.
EnvironmentFile={{ app_dir }}/shared/.env
EnvironmentFile=-{{ app_dir }}/shared/.env.release
ExecStart=/usr/bin/node {{ app_dir }}/current/apps/api/dist/main.js
Restart=on-failure
RestartSec=3
SyslogIdentifier=bestbooks-api

# Hardening (docs/05)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ app_dir }}/media
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
